# LowHttpServer.service
[Unit]
Description=Low-Tech Http Server
[Service]
Type=oneshot 
ExecStart=./LowHttpServer
[Install]
WantedBy=multi-user.target